# roles/supervisor/templates/mjq.ini
[program:meijiequan]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /usr/share/nginx/html/mjq/artisan queue:work --delay=10 --daemon
autostart=true
autorestart=true
user=root
numprocs=2
redirect_stderr=true
stdout_logfile=/opt/cron/logs/skyfire.log
